# AKSUMAEL/tools/systemd/aksumael-log-sync.timer
[Unit]
Description=Run aksumael-log-sync every 30s

[Timer]
OnBootSec=30s
OnUnitActiveSec=30s
AccuracySec=1s

[Install]
WantedBy=timers.target
